# etc/supervisord.conf
[supervisord]
directory = /home/jaakko/Repos/naprox
logfile = /home/jaakko/Repos/naprox/log/supervisord.log
pidfile = /home/jaakko/Repos/naprox/supervisord.pid
[program:naprox]
command = /home/jaakko/Repos/naprox/naprox.py
directory = /home/jaakko/Repos/naprox
logfile = /home/jaakko/Repos/naprox/log/supervisord.log
