# etc/supervisord.conf
[supervisord]
directory = /home/jaakko/Repos/naprox
logfile = /home/jaakko/Repos/naprox/log/supervisord.log
pidfile = /home/jaakko/Repos/naprox/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock

[inet_http_server]
port=127.0.0.1:8081

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock
username = unix
password = 1234
prompt = supervisord

[program:naprox]
command = /home/jaakko/Repos/naprox/naprox.py
directory = /home/jaakko/Repos/naprox
logfile = /home/jaakko/Repos/naprox/log/supervisord.log
